# data/init/drivnal.conf
description "drivnal daemon"

start on filesystem or runlevel [2345]
stop on runlevel [!2345]

respawn

exec drivnal --pidfile /var/run/drivnal.pid
